# edutrack.service
[Unit]
Description=CampusIQ School Management System
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/campusiq
# Next.js automatically loads .env.local from WorkingDirectory
# Do NOT use EnvironmentFile - it breaks on values with spaces
ExecStart=/usr/bin/node /home/ubuntu/campusiq/server.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=campusiq
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
